# hosts/macbook/com.hodgs004.screenlockmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hodgs004.screenlockmonitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Hodgs004/dotfiles/scripts/screen-lock-monitor.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/Hodgs004/.local/share/screen-lock-monitor.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Hodgs004/.local/share/screen-lock-monitor.err.log</string>
</dict>
</plist>
